# autostart/provider.desktop
[Desktop Entry]
Type=Application
Name=provider
Exec=/home/pi/run_provider.sh

